# packages/deb/debian.nboxui/postinst
#!/bin/sh -e

case "$1" in
    configure)
	# continue below
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
	exit 0
    ;;

    *)
	echo "postinst called with unknown argument \`$1'" >&2
	exit 0
    ;;
esac

umask 022

if ! grep -q nbox /etc/group; then
    echo 'Creating nbox group'
    /usr/sbin/groupadd -r nbox
fi

if ! /usr/bin/id -u nbox > /dev/null 2>&1; then
    echo "Creating nbox user..."
    /usr/sbin/useradd -M -N -g nbox -r nbox
fi

systemctl daemon-reload

systemctl enable nbox-scheduler.service > /dev/null
systemctl restart nbox-scheduler

systemctl enable nbox-rrd.service > /dev/null
systemctl restart nbox-rrd

# SNMPd
if [ -f /etc/snmp/snmpd.conf ]; then
    mv /etc/snmp/snmpd.conf /etc/snmp/snmpd.conf.bkp
    ln -s /etc/snmp/snmpd.conf.nbox /etc/snmp/snmpd.conf
    /etc/init.d/snmpd restart
fi
systemctl enable snmpd.service > /dev/null
systemctl restart snmpd

exit 0
