# packages/deb/debian.nboxui/prerm
#!/bin/sh -e

# Only shut the daemon down if we're really removing the package.  If this is
# an upgrade, we will instead do a restart in the postinst...

case "$1" in
upgrade)
	;;
*)
	systemctl disable nbox-rrd.service > /dev/null
	systemctl stop nbox-rrd

	systemctl disable nbox-scheduler.service > /dev/null
	systemctl stop nbox-scheduler
	;;
esac

exit 0
